# deploy/systemd/funding-scout-web.service
[Unit]
Description=funding-scout web UI (Dash + AG-Grid)
After=network-online.target postgresql.service funding-scout-snapshot.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/funding-scout
EnvironmentFile=/root/funding-scout/.env
# Host/port берутся из .env (FUNDING_SCOUT_WEB_HOST, FUNDING_SCOUT_WEB_PORT).
# Default 127.0.0.1 — loopback. Для LAN/VPN-доступа выставить интерфейс в .env.
ExecStart=/root/funding-scout/.venv/bin/funding-scout web --host ${FUNDING_SCOUT_WEB_HOST} --port ${FUNDING_SCOUT_WEB_PORT}
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

MemoryMax=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
